# tools/systemd/muxy_gunicorn.socket
[Unit]
Description=Muxy gunicorn socket

[Socket]
ListenStream=/run/muxy_gunicorn.sock

[Install]
WantedBy=sockets.target
