# tools/systemd/muxy-gunicorn.socket
[Unit]
Description=Muxy gunicorn socket

[Socket]
ListenStream=/run/muxy-gunicorn.sock

[Install]
WantedBy=sockets.target
